# server/config/echo.conf
[program:echo]
command=/home/stsouvallas/server-linux
autostart=true
autorestart=true
stderr_logfile=/var/log/echo-server.err.log
stdout_logfile=/var/log/echo-server.out.log